# beckend/tiktok-backend.service
[Unit]
Description=TikTok Scheduler - Backend API (FastAPI)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=magno
Group=magno
WorkingDirectory=/home/magno/work/tiktok-react/beckend
Environment="PATH=/home/magno/work/tiktok-react/beckend/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Comando para iniciar o backend
ExecStart=/home/magno/work/tiktok-react/beckend/venv/bin/python -m uvicorn src.http_health:app --host 0.0.0.0 --port 8082

# Restart automático
Restart=always
RestartSec=10

# Logs
StandardOutput=append:/home/magno/work/tiktok-react/beckend/logs/backend.log
StandardError=append:/home/magno/work/tiktok-react/beckend/logs/backend-error.log

# Limites
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
